# deploy/systemd/beckett-rpc.service
[Unit]
Description=Beckett Discord Rich Presence daemon
After=network.target graphical-session.target
# Keep trying even when Discord isn't open yet — but flag terminal failure like the main daemon.
StartLimitIntervalSec=0
OnFailure=beckett-alert@%N.service

[Service]
Type=simple
WorkingDirectory=%h/beckett
EnvironmentFile=-%h/.beckett/.env
Environment=DISCORD_RPC_LARGE_IMAGE=beckett
Environment=PATH=%h/.local/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env bun src/rpc/daemon.ts
Restart=always
RestartSec=15
# journald owns the logs (issue #30) — the old `append:~/.beckett/logs/rpc.log` grew unrotated.
StandardOutput=journal
StandardError=journal
# Unclean deaths alert (rate-limited hard for this unit — it's cosmetic and flaps by design).
ExecStopPost=%h/beckett/deploy/alert.sh %n

[Install]
WantedBy=default.target
