# deploy/systemd/beckett-v4.service
[Unit]
Description=Beckett v4 — multiplayer Plane ticket-queue (channel-scoped shared context; Opus concierge + poller + dispatcher)
After=network-online.target
Wants=network-online.target
# Crash-loop alerting (issue #30): 30 starts inside 10 minutes means something is wrong at boot —
# stop restarting (a tight loop burns CPU forever and helps nobody) and fire the terminal alert.
# Normal operation (a crash every so often + Restart=always) never gets near this limit.
StartLimitIntervalSec=600
StartLimitBurst=30
# %N (suffix stripped) — %n would double the suffix in the template instance name.
OnFailure=beckett-alert@%N.service

[Service]
Type=simple
WorkingDirectory=%h/beckett
ExecStart=/usr/local/bin/bun src/shell/v4-main.ts
# Always come back — a crash, an OOM kill, or a clean-looking exit all mean the coworker is
# offline until something restarts it. Crash recovery (issue #20) makes restarts safe.
Restart=always
RestartSec=3
# Pair with the dispatcher's graceful shutdown drain (20s budget): give SIGTERM time to commit
# WIP and persist the resume ledger before systemd escalates to SIGKILL.
TimeoutStopSec=30
# Secrets (DISCORD_TOKEN, PLANE_API_TOKEN, GITHUB_PAT, …) live in ~/.beckett/.env — loaded by
# the daemon itself too, but exporting here keeps child tooling consistent. `-` = optional.
EnvironmentFile=-%h/.beckett/.env
# The concierge + claude/pi/codex workers must resolve on PATH (a systemd unit does not inherit
# the interactive PATH). ~/.local/bin carries the modern node that pi needs (node >= 20).
# MUST stay in sync with `daemonPath()` in src/ops/doctor.ts — the doctor probes under this PATH.
Environment=PATH=%h/.local/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin
# Every unclean death posts a Discord alert within seconds (rate-limited in the script; clean
# stops are silent). Runs alongside Restart=always — the alert says systemd is already on it.
ExecStopPost=%h/beckett/deploy/alert.sh %n

[Install]
WantedBy=default.target
